# services.d/ble_bridge/run

#!/usr/bin/with-contenv bash
set -euo pipefail
# Delegate to the unified supervised entrypoint
exec /usr/src/app/run.sh
