# addon/services.d/echo_responder/run
#!/usr/bin/with-contenv bash
set -euo pipefail
OPT=/data/options.json
VENV="${VIRTUAL_ENV:-/opt/venv}"
cd /usr/src/app
exec "$VENV/bin/python" -m bb8_core.echo_responder
