# addon/services.d/ble_bridge/run
#!/usr/bin/with-contenv bash
set -euo pipefail

# ensure a writable log dir inside add-on data
mkdir -p /data/logs
export BB8_LOG_PATH="${BB8_LOG_PATH:-/data/logs/ha_bb8_addon.log}"

# Always start from /app so bb8_core is importable
cd /app

# PATH was set in Dockerfile to include /opt/venv/bin
exec /bin/bash /app/run.sh
